# RPI/files/overlay-root.service
[Unit]
Description=Boot with root filesystem overlay (read-only base)
After=multi-user.target
Before=cage@*.service

[Service]
Type=oneshot
ExecStart=/opt/scripts/overlay-root.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
