# RPI/files/overlay-root.service
[Unit]
Description=Boot with root filesystem overlay (read-only base)
DefaultDependencies=no
Before=sysinit.target local-fs.target
After=initrd.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/overlay-root.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
