# airootfs-shared/etc/systemd/user/set-pipewire-volume.service
[Unit]
Description=Monitor and set audio volume to 100%
After=pipewire.service pipewire-pulse.service
Wants=pipewire.service pipewire-pulse.service

[Service]
Type=simple
ExecStart=/opt/scripts/set-volume.sh
Restart=always
RestartSec=1

[Install]
WantedBy=default.target